# copy_to_system_etc/profile.d/xdg_cache_home.sh
if [ $USER ]; then
	export XDG_CACHE_HOME="/tmp/${USER}/.cache"
fi
